# contrib/kpatch.service
[Unit]
Description="Apply kpatch kernel patches"

[Service]
Type=oneshot
ExecStart=/usr/sbin/kpatch load --all

[Install]
WantedBy=multi-user.target
